# systemd-timers/mattermost-backup.service
[Unit]
Description=Mattermost Backup Service

[Service]
Type=oneshot
ExecStart=/path/to/backup.sh